# scripts/systemd/montagu-renew-certificate.service
[Service]
Type=oneshot
User=$USER
ExecStart=$MONTAGU renew-certificate $CONFIG
