# scripts/systemd/montagu-renew-certificate.service
[Service]
Type=oneshot
User=$USER
WorkingDirectory=$PWD
ExecStart=$MONTAGU renew-certificate
